# deploy/systemd/jack-code-server@.service
[Unit]
Description=Jack Code Server for %i
After=network.target

[Service]
Type=exec
User=%i
Environment=HOME=/home/%i
Environment=JACK_CODE_SERVER_TMPDIR=/tmp/jack-code-server-%i
ExecStart=/usr/bin/jack-code-server --bind-addr 127.0.0.1:8890 --user-data-dir /home/%i/.local/share/jack-code-server --extensions-dir /home/%i/.local/share/code-server/extensions
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
